# deploy/energyops-compose.service
[Unit]
Description=EnergyOps Docker Compose stack
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/var/www/energyops
ExecStart=/usr/bin/docker compose -f docker-compose.prod.yml up -d
ExecStop=/usr/bin/docker compose -f docker-compose.prod.yml down
TimeoutStartSec=600
Restart=on-failure

[Install]
WantedBy=multi-user.target
